# Unix/AIX/rhel-auditd-attack-deep.rules
# ===================================================================
# RHEL ATTACK DETECTION — DEEP RULESET
# File: rhel-auditd-attack-deep.rules
# Purpose: Deep detection of intrusion, privilege escalation, lateral
#          movement, credential access, persistence, and data theft.
# ===================================================================

##############################
# 1. RECEIVING Syscall Logs
##############################
-D
-b 8192
-f 1

##############################
# 2. ACCOUNT / AUTHENTICATION
##############################
-w /etc/passwd -p wa -k acct_change
-w /etc/shadow -p wa -k shadow_change
-w /etc/group -p wa -k group_change
-w /var/log/secure -p rwa -k auth_log

-a always,exit -F arch=b64 -S setuid -k setuid_calls
-a always,exit -F arch=b64 -S setgid -k setgid_calls

##############################
# 3. PRIVILEGE ESCALATION (EoP)
##############################
-w /usr/bin/sudo -p x -k sudo_exec
-a always,exit -F arch=b64 -S execve -C uid!=euid -k eop_exec
-a always,exit -F arch=b64 -S execve -C uid=0 -k root_command

##############################
# 4. COMMAND EXECUTION (Process Creation)
##############################
-a always,exit -F arch=b64 -S execve -k exec_log
-w /bin -p x -k bin_exec
-w /usr/bin -p x -k usrbin_exec
-w /usr/sbin -p x -k usrsbin_exec

##############################
# 5. NETWORK & REMOTE SERVICES
##############################
-w /etc/ssh/sshd_config -p wa -k ssh_config
-w /etc/hosts.allow -p wa -k hosts_allow
-w /etc/hosts.deny -p wa -k hosts_deny

-a always,exit -F arch=b64 -S connect -k net_connect
-a always,exit -F arch=b64 -S accept  -k net_accept

##############################
# 6. FILE ACCESS & TAMPERING
##############################
-w /etc/sudoers -p wa -k sudoers_change
-w /etc/cron.allow -p wa -k cron_allow
-w /etc/cron.deny -p wa -k cron_deny
-w /etc/cron.d -p wa -k cron_d
-w /etc/crontab -p wa -k crontab

-w /etc/systemd/system -p wa -k systemd_service
-w /lib/systemd/system -p wa -k systemd_lib

##############################
# 7. PERSISTENCE MECHANISMS
##############################
-w /etc/rc.local -p wa -k rc_local
-w /etc/init.d -p wa -k init_scripts
-w /etc/profile -p wa -k bash_profile
-w /root/.bashrc -p wa -k root_bashrc
-w /home/*/.ssh -p wa -k ssh_keys

##############################
# 8. CREDENTIAL ACCESS
##############################
-w /etc/ssh/sshd_config -p wa -k ssh_config
-w /root/.ssh -p wa -k root_ssh_keys
-w /home/*/.ssh -p wa -k user_ssh_keys

##############################
# 9. KERNEL / SECURITY SETTINGS
##############################
-w /etc/sysctl.conf -p wa -k sysctl_change
-w /etc/audit -p wa -k audit_conf

-a always,exit -F arch=b64 -S init_module -k kernel_module
-a always,exit -F arch=b64 -S delete_module -k kernel_module

##############################
# 10. SUSPICIOUS BINARY EXECUTION
##############################
-a always,exit -F arch=b64 -S execve -F path=/tmp -k tmp_exec
-a always,exit -F arch=b64 -S execve -F dir=/dev/shm -k shm_exec
-a always,exit -F arch=b64 -S chmod -k suspicious_perm
-a always,exit -F arch=b64 -S chown -k suspicious_owner

##############################
# 11. DATA EXFILTRATION (File Read)
##############################
-a always,exit -F arch=b64 -S open -F dir=/etc -k file_read_etc
-a always,exit -F arch=b64 -S open -F dir=/var/www -k web_leak
-a always,exit -F arch=b64 -S open -F dir=/home -k home_read

##############################
# 12. CLEANUP / COVER TRACKS
##############################
-w /var/log -p wa -k log_tamper
-w /usr/bin/wipe -p x -k wipe_exec
-w /usr/bin/shred -p x -k shred_exec

##############################
# END OF FILE
##############################
