# Ubuntu/auditd/ubuntu-auditd-deep-plura-v3-latest..rules
# PLURA_CONFIG_VERSION: ubuntu-auditd-deep-plura-v3-latest..rules
# PLURA_SCHEMA_TARGET: v3.1
# Profile: DEEP (incident response / higher visibility)
# Notes:
# - Includes BASELINE rules + deeper syscall auditing.
# - Default deep focuses on authenticated sessions (auid!=unset) and root actions (euid=0) to control volume.
# - Optional "all user actions" blocks are provided but commented.
# Ubuntu adaptation:
# - Baseline block is inlined with Ubuntu path adjustments.
# - Deep syscall rules are OS-agnostic and kept identical to RHEL v3.1.

############################################################
# 0) Include BASELINE content
#    (For repo management: keep the baseline block identical, then append deep block.)
############################################################

# --- START: baseline v3.1 rules (Ubuntu adapted) ---
# PLURA_CONFIG_VERSION: auditd-baseline-plura-ubuntu-v3.1.rules
# PLURA_SCHEMA_TARGET: v3.1
# Profile: BASELINE (low-noise)
# Notes:
# - Removed high-volume directory execute watches (/run, /var/run, /var/lock, /dev/shm).
# - Keep legacy qubit-* keys where reasonable, and add plura_* keys for coverage.
# - Avoid watching paths that may not exist by default (kept as optional/commented).
# Ubuntu adaptation:
# - Replace SELinux path watch with AppArmor paths (keep key qubit-selinux for backward compatibility).
# - Replace systemd unit dir /usr/lib/systemd/system with /lib/systemd/system (Ubuntu default).
# - Replace /etc/bashrc with /etc/bash.bashrc.
# - Replace RHEL iptables config path (/etc/sysconfig/iptables) with Ubuntu firewall options (UFW / nftables / iptables-persistent) as optional.

############################################################
# 1) Legacy (qubit-*) - trimmed for noise
############################################################

# Identity files (high signal)
-w /etc/shadow -p wa -k qubit-shadow
-w /etc/passwd -p wa -k qubit-passwd

# Privilege/permission tooling (moderate signal)
-w /usr/bin/chown -p x -k qubit-chown
-w /usr/bin/chmod -p x -k qubit-perm

# Compiler execution
-w /usr/bin/gcc -p x -k qubit-gcc

# Recon / account tooling
-w /usr/bin/whoami -p x -k qubit-recon
-w /usr/bin/w -p x -k qubit-recon
-w /usr/sbin/useradd -p x -k qubit-useradd
-w /usr/sbin/userdel -p x -k qubit-userdel

# Network / transfer tools
-w /usr/bin/ncat -p x -k qubit-nc
-w /usr/bin/nc -p x -k qubit-nc
-w /usr/bin/curl -p x -k qubit-curl
-w /usr/bin/wget -p x -k qubit-wget

# Utility / core config
-w /usr/bin/base64 -p x -k qubit-base64
-w /etc/sysctl.conf -p wa -k qubit-kernel
# Ubuntu: AppArmor config (replaces /etc/selinux) - keep legacy key for compatibility
-w /etc/apparmor -p wa -k qubit-selinux
-w /etc/apparmor.d -p wa -k qubit-selinux
# Optional (enable if present)
# -w /etc/default/apparmor -p wa -k qubit-selinux
-w /etc/hosts -p wa -k qubit-host
-w /etc/localtime -p wa -k qubit-time

# Remote access (as-is)
-w /usr/bin/vncserver -p x -k qubit-vnc
-w /usr/bin/x11vnc -p x -k qubit-vnc
-w /usr/bin/Xvnc -p x -k qubit-vnc
-w /usr/bin/sftp -p x -k qubit-sftp
-w /usr/bin/rsync -p x -k qubit-rsync
-w /usr/bin/scp -p x -k qubit-scp

# IMPORTANT: Removed (high-volume) directory execute watches from v3.0:
# -w /run -p x -k qubit-dir-/run
# -w /var/run -p x -k qubit-dir-/var/run
# -w /var/lock -p x -k qubit-dir-/var/lock
# -w /dev/shm -p x -k qubit-dir-/dev/shm
#
# Removed (often high-frequency in ops environments):
# -w /usr/bin/rm -p x -k qubit-rm
# -w /usr/bin/cp -p x -k qubit-cp

# Socket syscall (AF_PACKET=17) - add auid filter to reduce daemon noise
-a always,exit -F arch=b64 -S socket -F a0=17 -F auid!=unset -k qubit-socket-17
-a always,exit -F arch=b32 -S socket -F a0=17 -F auid!=unset -k qubit-socket-17

# Firewall (Ubuntu)
# iptables binaries (may be nft-backed depending on system alternatives)
-w /usr/sbin/iptables -p x -k qubit-iptables
-w /sbin/iptables -p x -k qubit-iptables
# Optional (enable if present)
# -w /usr/sbin/ip6tables -p x -k qubit-iptables
# -w /sbin/ip6tables -p x -k qubit-iptables
# -w /usr/sbin/nft -p x -k qubit-iptables
# -w /etc/nftables.conf -p wa -k qubit-iptables
# Ubuntu default firewall (UFW)
# -w /usr/sbin/ufw -p x -k qubit-iptables
# -w /etc/ufw -p wa -k qubit-iptables
# -w /etc/default/ufw -p wa -k qubit-iptables
# iptables-persistent (if used)
# -w /etc/iptables -p wa -k qubit-iptables
# -w /etc/iptables/rules.v4 -p wa -k qubit-iptables
# -w /etc/iptables/rules.v6 -p wa -k qubit-iptables


############################################################
# 2) PLURA baseline coverage (plura_*)
############################################################

# Identity / auth db (add missing group/gshadow)
-w /etc/group -p wa -k plura_identity
-w /etc/gshadow -p wa -k plura_identity

# User/group management commands
-w /usr/sbin/usermod -p x -k plura_user_mgmt
-w /usr/sbin/groupadd -p x -k plura_group_mgmt
-w /usr/sbin/groupdel -p x -k plura_group_mgmt
-w /usr/sbin/groupmod -p x -k plura_group_mgmt

# Password/account attribute changes
-w /usr/bin/passwd -p x -k plura_passwd_change
-w /usr/bin/chage -p x -k plura_passwd_change
-w /usr/sbin/chpasswd -p x -k plura_passwd_change

# Privilege escalation & sudo config
-w /bin/su -p x -k plura_priv_exec
-w /usr/bin/sudo -p x -k plura_priv_exec
-w /usr/bin/sudoedit -p x -k plura_priv_exec
-w /etc/sudoers -p wa -k plura_sudo_config
-w /etc/sudoers.d -p wa -k plura_sudo_config

# SSH server config changes
-w /etc/ssh/sshd_config -p wa -k plura_ssh_config
-w /etc/ssh -p wa -k plura_ssh_config

# Persistence - systemd and init scripts
-w /etc/systemd/system -p wa -k plura_systemd
-w /lib/systemd/system -p wa -k plura_systemd
# Optional (enable if present / usr-merge environments)
# -w /usr/lib/systemd/system -p wa -k plura_systemd
-w /etc/init.d -p wa -k plura_init

# Persistence - cron/anacron
-w /etc/crontab -p wa -k plura_cron
-w /etc/cron.d -p wa -k plura_cron
-w /etc/cron.hourly -p wa -k plura_cron
-w /etc/cron.daily -p wa -k plura_cron
-w /etc/cron.weekly -p wa -k plura_cron
-w /etc/cron.monthly -p wa -k plura_cron
# Debian/Ubuntu cron spool location
-w /var/spool/cron/crontabs -p wa -k plura_cron
# -w /var/spool/cron -p wa -k plura_cron
-w /etc/anacrontab -p wa -k plura_cron

# Auditd tampering (Defense Evasion)
-w /etc/audit -p wa -k plura_audit_config
-w /etc/audit/rules.d -p wa -k plura_audit_config
-w /etc/audit/auditd.conf -p wa -k plura_audit_config
-w /usr/sbin/auditctl -p x -k plura_audit_tools
-w /usr/sbin/augenrules -p x -k plura_audit_tools
-w /usr/sbin/auditd -p x -k plura_audit_tools

# PAM/auth configuration tampering
-w /etc/pam.d -p wa -k plura_auth_config
-w /etc/security -p wa -k plura_auth_config
-w /etc/login.defs -p wa -k plura_auth_config

# Loader / library path hijacking (ld.so.preload is often absent -> optional)
-w /etc/ld.so.conf -p wa -k plura_loader
-w /etc/ld.so.conf.d -p wa -k plura_loader
# -w /etc/ld.so.preload -p wa -k plura_loader   # enable only if file exists

# Shell/profile persistence
-w /etc/profile -p wa -k plura_shell_profile
-w /etc/profile.d -p wa -k plura_shell_profile
-w /etc/bash.bashrc -p wa -k plura_shell_profile
# Optional (for distros that use /etc/bashrc)
# -w /etc/bashrc -p wa -k plura_shell_profile

# Logging pipeline tampering
-w /etc/rsyslog.conf -p wa -k plura_logging
-w /etc/rsyslog.d -p wa -k plura_logging
-w /etc/logrotate.conf -p wa -k plura_logging
-w /etc/logrotate.d -p wa -k plura_logging

# DNS/network basics (usually exist)
-w /etc/resolv.conf -p wa -k plura_net_conf
-w /etc/hostname -p wa -k plura_net_conf

# High-signal syscall rules (rare events)
# Kernel module load/unload by authenticated session (avoid daemon noise)
-a always,exit -F arch=b64 -S init_module -S delete_module -F euid=0 -F auid!=unset -k plura_kernel_module
-a always,exit -F arch=b32 -S init_module -S delete_module -F euid=0 -F auid!=unset -k plura_kernel_module
# Optional if supported
# -a always,exit -F arch=b64 -S finit_module -F euid=0 -F auid!=unset -k plura_kernel_module
# -a always,exit -F arch=b32 -S finit_module -F euid=0 -F auid!=unset -k plura_kernel_module

# System time changes (rare / high-signal)
-a always,exit -F arch=b64 -S adjtimex -S settimeofday -S clock_settime -F euid=0 -F auid!=unset -k plura_time_change
-a always,exit -F arch=b32 -S adjtimex -S settimeofday -S clock_settime -F euid=0 -F auid!=unset -k plura_time_change
# --- END: baseline v3.1 rules ---


############################################################
# 3) PLURA deep additions (syscall rules)
############################################################

# 3-1) Root command execution from authenticated sessions (very useful, manageable)
-a always,exit -F arch=b64 -S execve -F euid=0 -F auid!=unset -k plura_execve_root
-a always,exit -F arch=b32 -S execve -F euid=0 -F auid!=unset -k plura_execve_root
# Optional if supported by kernel:
# -a always,exit -F arch=b64 -S execveat -F euid=0 -F auid!=unset -k plura_execve_root
# -a always,exit -F arch=b32 -S execveat -F euid=0 -F auid!=unset -k plura_execve_root

# 3-2) Privilege/identity changes (set*id)
-a always,exit -F arch=b64 -S setuid -S setreuid -S setresuid -S setgid -S setregid -S setresgid -F auid!=unset -k plura_sys_privchange
-a always,exit -F arch=b32 -S setuid -S setreuid -S setresuid -S setgid -S setregid -S setresgid -F auid!=unset -k plura_sys_privchange

# 3-3) Capability changes
-a always,exit -F arch=b64 -S capset -F auid!=unset -k plura_sys_capset
-a always,exit -F arch=b32 -S capset -F auid!=unset -k plura_sys_capset

# 3-4) Mount/unmount (persistence / defense evasion)
-a always,exit -F arch=b64 -S mount -S umount2 -F euid=0 -F auid!=unset -k plura_fs_mount
-a always,exit -F arch=b32 -S mount -S umount2 -F euid=0 -F auid!=unset -k plura_fs_mount

# 3-5) File tampering (delete/rename/rmdir) - default: root only
-a always,exit -F arch=b64 -S unlink -S unlinkat -S rename -S renameat -S rmdir -F euid=0 -F auid!=unset -k plura_file_tamper_root
-a always,exit -F arch=b32 -S unlink -S unlinkat -S rename -S renameat -S rmdir -F euid=0 -F auid!=unset -k plura_file_tamper_root

# Optional: all authenticated users (can be noisy depending on workload)
# -a always,exit -F arch=b64 -S unlink -S unlinkat -S rename -S renameat -S rmdir -F auid!=unset -k plura_file_tamper_user
# -a always,exit -F arch=b32 -S unlink -S unlinkat -S rename -S renameat -S rmdir -F auid!=unset -k plura_file_tamper_user

# 3-6) Permission/ownership changes - default: root only
-a always,exit -F arch=b64 -S chmod -S fchmod -S fchmodat -S chown -S fchown -S fchownat -S lchown -F euid=0 -F auid!=unset -k plura_file_permchg_root
-a always,exit -F arch=b32 -S chmod -S fchmod -S fchmodat -S chown -S fchown -S fchownat -S lchown -F euid=0 -F auid!=unset -k plura_file_permchg_root

# Optional: all authenticated users
# -a always,exit -F arch=b64 -S chmod -S fchmod -S fchmodat -S chown -S fchown -S fchownat -S lchown -F auid!=unset -k plura_file_permchg_user
# -a always,exit -F arch=b32 -S chmod -S fchmod -S fchmodat -S chown -S fchown -S fchownat -S lchown -F auid!=unset -k plura_file_permchg_user

# 3-7) xattr changes (ACL/hidden attributes) - default: root only
-a always,exit -F arch=b64 -S setxattr -S lsetxattr -S fsetxattr -S removexattr -S lremovexattr -S fremovexattr -F euid=0 -F auid!=unset -k plura_xattr_root
-a always,exit -F arch=b32 -S setxattr -S lsetxattr -S fsetxattr -S removexattr -S lremovexattr -S fremovexattr -F euid=0 -F auid!=unset -k plura_xattr_root

# 3-8) ptrace (process injection / debugging abuse)
-a always,exit -F arch=b64 -S ptrace -F auid!=unset -k plura_ptrace
-a always,exit -F arch=b32 -S ptrace -F auid!=unset -k plura_ptrace

# Optional (advanced / kernel-dependent): memfd_create is not available on older kernels
# -a always,exit -F arch=b64 -S memfd_create -F auid!=unset -k plura_memfd
# -a always,exit -F arch=b32 -S memfd_create -F auid!=unset -k plura_memfd
