# RHEL/auditd/auditd-baseline-plura.rules
# PLURA_CONFIG_VERSION: auditd-baseline-plura.rules
# PLURA_SCHEMA_TARGET: v3.0

###############################################
# 1) Existing (qubit-*)
###############################################
# 비밀번호 정보 변경 (NO 4)
-w /etc/shadow -p wa -k qubit-shadow
# 소유권 및 권한 변경 도구 (NO 159, 158)
-w /usr/bin/chown -p x -k qubit-chown
-w /usr/bin/chmod -p x -k qubit-perm
# 컴파일러 실행 (NO 157)
-w /usr/bin/gcc -p x -k qubit-gcc
# 계정 정보 및 관리 (NO 156, 51, 49, 8, 5)
-w /etc/passwd -p wa -k qubit-passwd
-w /usr/bin/whoami -p x -k qubit-recon
-w /usr/sbin/userdel -p x -k qubit-userdel
-w /usr/bin/w -p x -k qubit-recon
-w /usr/sbin/useradd -p x -k qubit-useradd
# 네트워크 및 전송 도구 (NO 29, 28, 26, 25)
-w /usr/bin/ncat -p x -k qubit-nc
-w /usr/bin/nc -p x -k qubit-nc
-w /usr/bin/curl -p x -k qubit-curl
-w /usr/bin/wget -p x -k qubit-wget
# 시스템 유틸리티 및 설정 파일 (NO 27, 12, 11, 10, 9)
-w /usr/bin/base64 -p x -k qubit-base64
-w /etc/sysctl.conf -p wa -k qubit-kernel
-w /etc/selinux/ -p wa -k qubit-selinux
-w /etc/hosts -p wa -k qubit-host
-w /etc/localtime -p wa -k qubit-time
# VNC 및 원격 접속 (NO 213, 212, 211, 210, 208)
-w /usr/bin/vncserver -p x -k qubit-vnc
-w /usr/bin/x11vnc -p x -k qubit-vnc
-w /usr/bin/Xvnc -p x -k qubit-vnc
-w /usr/bin/sftp -p x -k qubit-sftp
-w /usr/bin/rsync -p x -k qubit-rsync
# 런타임 및 Socket 호출 (NO 192, 191, 190)
-w /run -p x -k qubit-dir-/run
-a always,exit -F arch=b32 -S socket -F a0=17 -k qubit-socket-17
-a always,exit -F arch=b64 -S socket -F a0=17 -k qubit-socket-17
# 파일 조작 명령어 (NO 189, 188, 187)
-w /usr/bin/scp -p x -k qubit-scp
-w /usr/bin/rm -p x -k qubit-rm
-w /usr/bin/cp -p x -k qubit-cp
# IPTables 방화벽 (NO 186, 185, 184, 183)
-w /etc/iptables -p wa -k qubit-iptables
-w /usr/sbin/iptables -p x -k qubit-iptables
-w /etc/sysconfig/iptables -p wa -k qubit-iptables
-w /sbin/iptables -p x -k qubit-iptables
# 주요 디렉토리 감시 (NO 182, 181, 180)
-w /var/run -p x -k qubit-dir-/var/run
-w /var/lock -p x -k qubit-dir-/var/lock
-w /dev/shm -p x -k qubit-dir-/dev/shm

###############################################
# 2) Added (plura attack baseline) - duplicates skipped
###############################################
###############################################
# rhel-auditd-attack-baseline.rules
# - 목적: MITRE ATT&CK 최소 커버용 PLURA 기본 룰
# - 특징: auid>=1000 (실사용자) 중심, 운영 가능한 로그량
# - 태그: plura_* 로 시작 (PLURA 에이전트 생성 룰 식별용)
###############################################

##############
# 0. 공통 필터
##############
# (여기서는 -D, -b 등 글로벌 옵션은 다루지 않음)
# 기본 정책은 상위 메인 rules 파일에서 관리한다고 가정.

#######################################
# 1. 계정 / 그룹 / 인증 정보 변경 추적
#######################################

# 계정/그룹 데이터베이스 변경 (T1136 Account Creation, T1098 Account Manipulation)
# (skip duplicate watch) -w /etc/passwd   -p wa -k plura_identity
# (skip duplicate watch) -w /etc/shadow   -p wa -k plura_identity
-w /etc/group    -p wa -k plura_identity
-w /etc/gshadow  -p wa -k plura_identity

# 계정/그룹 관리 명령 사용 추적
# (skip duplicate watch) -w /usr/sbin/useradd   -p x -k plura_user_mgmt
# (skip duplicate watch) -w /usr/sbin/userdel   -p x -k plura_user_mgmt
-w /usr/sbin/usermod   -p x -k plura_user_mgmt
-w /usr/sbin/groupadd  -p x -k plura_group_mgmt
-w /usr/sbin/groupdel  -p x -k plura_group_mgmt
-w /usr/sbin/groupmod  -p x -k plura_group_mgmt

# 패스워드/계정 속성 변경
-w /usr/bin/passwd     -p x -k plura_passwd_change
-w /usr/bin/chage      -p x -k plura_passwd_change
-w /usr/sbin/chpasswd  -p x -k plura_passwd_change

##########################################
# 2. 로그인/세션/권한상승(계정 탈취, sudo/su)
##########################################

# su, sudo 실행 (권한 상승 시도 추적 – T1078 Valid Accounts, T1068 Priv Esc)
-w /bin/su         -p x -k plura_su_exec
-w /usr/bin/sudo   -p x -k plura_sudo_exec
-w /usr/bin/sudoedit -p x -k plura_sudo_exec

# sudoers 설정 변경 (Defense Evasion)
-w /etc/sudoers       -p wa -k plura_sudo_config
-w /etc/sudoers.d/    -p wa -k plura_sudo_config

##########################################
# 3. 프로세스/명령 실행 (EXECVE – ATT&CK 핵심)
##########################################

# 일반 사용자(auid>=1000)의 execve 전체 추적 (T1059 Command and Scripting)
# 64bit
-a always,exit -F arch=b64 -S execve \
  -C auid>=1000 -F auid!=4294967295 \
  -k plura_execve_user

# 32bit (멀티아치 시스템일 경우만 활성화) 
#-a always,exit -F arch=b32 -S execve \
#  -C auid>=1000 -F auid!=4294967295 \
#  -k plura_execve_user

##########################################
# 4. 중요 설정 파일 변경 (Persistence / Defense Evasion)
##########################################

# SSH 서버 설정 변경
-w /etc/ssh/sshd_config   -p wa -k plura_ssh_config

# systemd 서비스 유닛 변경 (서비스 등록/변조 – Persistence)
-w /etc/systemd/system/       -p wa -k plura_service_change
-w /usr/lib/systemd/system/   -p wa -k plura_service_change

# Crontab 및 주기 작업 (T1053 Scheduled Task/Job)
-w /etc/crontab               -p wa -k plura_cron
-w /etc/cron.hourly/          -p wa -k plura_cron
-w /etc/cron.daily/           -p wa -k plura_cron
-w /etc/cron.weekly/          -p wa -k plura_cron
-w /etc/cron.monthly/         -p wa -k plura_cron
-w /etc/cron.d/               -p wa -k plura_cron
-w /var/spool/cron/           -p wa -k plura_cron

##########################################
# 5. 방화벽 / 네트워크 정책 변경 (Netfilter)
##########################################

# iptables / nftables 실행 추적 (T1562 Defense Evasion, 방화벽 변조)
# (skip duplicate watch) -w /usr/sbin/iptables   -p x -k plura_fw_change
-w /usr/sbin/ip6tables  -p x -k plura_fw_change
-w /usr/sbin/nft        -p x -k plura_fw_change

##########################################
# 6. SELinux / MAC 정책 변경
##########################################

# SELinux 설정/정책 파일 변경
# (skip duplicate watch) -w /etc/selinux/            -p wa -k plura_selinux_conf

##########################################
# 7. Auditd / 로깅 설정 변경
##########################################

# auditd 설정 자체 변경 탐지 (Defense Evasion – 로그 끄기)
-w /etc/audit/           -p wa -k plura_audit_config
-w /usr/sbin/auditctl    -p x  -k plura_audit_config
-w /usr/sbin/auditd      -p x  -k plura_audit_config

##########################################
# 8. 핵심 바이너리 위변조 (무결성 관점 최소)
##########################################

# su / sudo / sshd / passwd 등 핵심 바이너리 변경
-w /bin/su               -p wa -k plura_core_bin
-w /usr/bin/sudo         -p wa -k plura_core_bin
-w /usr/sbin/sshd        -p wa -k plura_core_bin
-w /usr/bin/passwd       -p wa -k plura_core_bin
